# crontab
*  *  *  *  *  /metrics/metrics.sh > /dev/stdout